# deploy/dexity-server.service
[Unit]
Description=Dexity backend (Fastify)
After=network.target

[Service]
Type=simple
User=dexity
WorkingDirectory=/srv/dexity/current/server
EnvironmentFile=/srv/dexity/.env
ExecStart=/usr/bin/node dist/server/src/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
